# wireguard-ui/rootfs/etc/services.d/refresh-wg/run
#!/usr/bin/with-contenv bashio
# shellcheck disable=SC2191

set -euo pipefail

# shellcheck disable=SC1091
source /usr/local/share/service-functions.sh

bashio::log.info 'Starting wg0.conf monitor...'

dmz_subnet="$(try_config_or_default 'dmz_subnet | join(" ")?')"
dmz_subnet_allow="$(try_config_or_default 'dmz_subnet_allow | join(" ")?')"
allowance_rule_separator='|'
export allowance_rule_separator dmz_subnet dmz_subnet_allow
if [ "${dmz_subnet:-}" = null ]; then
  unset dmz_subnet
fi

user_loglevel=$(try_config_or_default 'log_level')
if [ "$user_loglevel" = trace ]; then
  set -x
fi

default_nic="$(ip route | grep default | awk '{print $5}')"
default_route_ip="$(ip route | awk '$1 == "default" { print $3 }')"
export default_nic default_route_ip

# daemon restarts wg interface on modification
# [Interface] の Address 行から、カンマ区切り(例: "192.168.101.0/24, fd42:42:42::0/64")
# の各アドレスを 1 行ずつ取り出す。IPv4 と IPv6 が混在していても両方拾えるようにする。
iface_addr() {
  awk -F'=' '
    /^[[:space:]]*\[/ { iface = 0 }
    /^[[:space:]]*\[Interface\][[:space:]]*$/ { iface = 1; next }
    iface && $1 ~ /^[[:space:]]*Address[[:space:]]*$/ {
      val = $2
      gsub(/[[:space:]]/, "", val)
      n = split(val, a, ",")
      for (i = 1; i <= n; i++) if (a[i] != "") print a[i]
    }
  ' /etc/wireguard/wg0.conf
}
is_ipv6() { case "$1" in *:*) return 0 ;; *) return 1 ;; esac; }
clear_chain_rule() {
  case "$#" in
    1)
      iptables-save | \
        grep -F -- "-A ${1}" | \
        sed 's/-A/-D/' | \
        tr '\n' '\0' | \
        xargs -0 -I{} /bin/bash -c 'iptables {}'
      ;;
    2)
      iptables-save | \
        grep -F -- "-A ${1}" | \
        grep -F -- "-j ${2}" | \
        sed 's/-A/-D/' | \
        tr '\n' '\0' | \
        xargs -0 -I{} /bin/bash -c 'iptables {}'
      ;;
  esac
}
clear_subnet_chain_rules() (
  set +o pipefail
  clear_chain_rule OUTPUT WG_DMZ
  clear_chain_rule FORWARD WG_DMZ
  clear_chain_rule WG_DMZ_allow
)
add_chain_rules() {
  # 1 - dmz_subnet
  # 2 - dmz_subnet_allow
  # 3 - WG_DMZ
  # 4 - WG_DMZ_allow
  echo "$1" | xargs -n1 | xargs -I{} -- iptables -I OUTPUT -s {} -j "$3"
  echo "$1" | xargs -n1 | xargs -I{} -- iptables -I FORWARD -s {} -j "$3"
  while read -r allow_rule; do
    apply_allow_rule "$allow_rule" "$1" "$4"
  done <<< "$(xargs -n1 <<< "${2:-}")"
}
clear_nats() (
  set +o pipefail
  if iptables-save 2>/dev/null | grep -F -- '-j MASQUERADE' > /dev/null; then
    iptables-save | grep -F -- '-j MASQUERADE' | sed 's/^-A/-D/' | tr '\n' '\0' | xargs -0 -I{} /bin/bash -c 'iptables -t nat {}'
  fi
  # IPv6 側の MASQUERADE も同様に掃除する(ip6tables が無い/使えない環境では無視)
  if ip6tables-save 2>/dev/null | grep -F -- '-j MASQUERADE' > /dev/null; then
    ip6tables-save | grep -F -- '-j MASQUERADE' | sed 's/^-A/-D/' | tr '\n' '\0' | xargs -0 -I{} /bin/bash -c 'ip6tables -t nat {}'
  fi
)
# Address のうち IPv4 は iptables、IPv6 は ip6tables で MASQUERADE を張る。
# IPv6 NAT に未対応の環境では警告を出すだけに留め、IPv4 の動作は止めない。
nat_ip_ranges() {
  local addr
  iface_addr | while read -r addr; do
    [ -z "${addr}" ] && continue
    if is_ipv6 "${addr}"; then
      if ! ip6tables -t nat -A POSTROUTING -s "${addr}" -o "${default_nic}" -j MASQUERADE 2>/dev/null; then
        echo "Warning: could not add IPv6 MASQUERADE for ${addr} (IPv6 NAT unsupported on this host?)"
      fi
    else
      iptables -t nat -A POSTROUTING -s "${addr}" -o "${default_nic}" -j MASQUERADE
    fi
  done
}
checksum_interface() {
  if [ ! -f /etc/wireguard/wg0.conf ]; then
    return
  fi
  awk '$1 ~ /^\[Interface\]$/ { print $1; out=1; next }; out != 1 { next; }; $1 ~ /^[^\[#]/ { print; next}; { exit }' /etc/wireguard/wg0.conf | \
    sha256sum
}
refresh_wireguard() {
  local current_checksum
  current_checksum="$(checksum_interface)"
  if [ "${old_interface_checksum}" = "$current_checksum" ]; then
    wg syncconf wg0 <(wg-quick strip wg0)
    echo "Wireguard clients synchronized preserving wg0 interface."
  else
    echo "Interface changed... recreating wg0"
    old_interface_checksum="$current_checksum"
    if ip link show wg0; then
      wg-quick down wg0
    fi
    wg-quick up wg0
  fi
  clear_nats
  clear_subnet_chain_rules
  if [ -n "${dmz_subnet:-}" ]; then
    add_chain_rules "${dmz_subnet}" "${dmz_subnet_allow:-}" WG_DMZ WG_DMZ_allow
  fi
  nat_ip_ranges
  if [ -n "${dmz_subnet:-}" ]; then
    iptables -nL
  fi
}
# DMZ exclusion functions
is_cidr_or_ip() {
  local val
  if [ "$#" = 0 ]; then
    val="$(cat)"
  else
    val="$1"
  fi
  if [ "${val:-}" = default_route_ip ]; then
    return 0
  fi
  # stock alpine image: validate with grep instead of python3
  grep -Eq '^([0-9]{1,3}\.){3}[0-9]{1,3}(/[1-3]?[0-9])?$' <<< "${val:-}"
}
get_src_net() {
  case "$(count_sep "$1")" in
    0)
      return
      ;;
    1)
      if ! {
          grep '|$' > /dev/null <<< "$1" || {
            cut -d'|' -f1 <<< "$1" | is_cidr_or_ip && \
            cut -d'|' -f2 <<< "$1" | is_cidr_or_ip
          }
        }; then
        return
      fi
      ;;
  esac
  cut -d'|' -f1 <<< "$1"
}
get_dst_net() {
  case "$(count_sep "$1")" in
    0)
      if is_cidr_or_ip "$1"; then
        echo "$1"
      fi
      ;;
    1)
      if cut -d'|' -f2 <<< "$1" | is_cidr_or_ip; then
        cut -d'|' -f2 <<< "$1"
      elif grep '|$' > /dev/null <<< "$1"; then
        return
      else
        cut -d'|' -f1 <<< "$1"
      fi
      ;;
    2)
      cut -d'|' -f2 <<< "$1"
      ;;
  esac
}
get_port_field() {
  local port
  case "$(count_sep "$1")" in
    0)
      if ! grep '^[:,0-9]*/[a-z]*$' > /dev/null <<< "$1"; then
        return
      fi
      port="$1"
      ;;
    1)
      if {
          cut -d'|' -f1 <<< "$1" | is_cidr_or_ip && \
          cut -d'|' -f2 <<< "$1" | is_cidr_or_ip
        }; then
        return
      fi
      port="$(cut -d'|' -f2 <<< "$1")"
      ;;
    2)
      port="$(cut -d'|' -f3 <<< "$1")"
      ;;
  esac
  echo "$port"
}
get_dst_port() {
  local port
  port="$(get_port_field "$1")"
  echo "${port%/*}"
}
get_dst_proto() {
  local port
  port="$(get_port_field "$1")"
  echo "${port#*/}"
}
count_sep() {
  echo "$1" | grep -o "${allowance_rule_separator}" | tr -d '\n' | wc -c
}
add_exclusion_rule_to_iptables() (
  iptables_args=( iptables -A "${1}" )
  if [ -n "${src_net:-}" ]; then
    iptables_args+=( -s "${src_net}" )
  fi
  if [ -n "${dst_net:-}" ]; then
    iptables_args+=( -d "${dst_net}" )
  fi
  if [ "${dst_proto:-}" = icmp ]; then
    iptables_args+=( -p icmp )
  elif [ -n "${dst_port:-}" ] && [ -n "${dst_proto:-}" ]; then
    case "${dst_proto}" in
      tcp|udp)
        iptables_args+=( -p "${dst_proto}" )
        ;;
    esac
    if grep '^[0-9]*$' > /dev/null <<< "$dst_port"; then
      iptables_args+=( -m "${dst_proto}" --dport "${dst_port}" )
    else
      iptables_args+=( -m multiport --dports "${dst_port}" )
    fi
  fi

  if [ "${dst_proto:-}" = icmp ]; then
    set -x
    # all ICMP traffic is recognized as RELATED except for pings
    "${iptables_args[@]}" --icmp-type echo-request -j ACCEPT
  else
    set -x
    "${iptables_args[@]}" -j ACCEPT
  fi
)
set_parsed_environment_variables() {
  src_net="$(get_src_net "$1")"
  dst_net="$(get_dst_net "$1")"
  dst_port="$(get_dst_port "$1")"
  dst_proto="$(get_dst_proto "$1")"
  export src_net dst_net dst_port dst_proto
  if [ "${src_net:-}" = default_route_ip ]; then
    src_net="${default_route_ip}"
  fi
  if [ "${dst_net:-}" = default_route_ip ]; then
    dst_net="${default_route_ip}"
  fi
}
apply_allow_rule() (
  # 1 - allow_rule
  # 2 - dmz_subnet
  # 3 - WG_DMZ_allow
  set_parsed_environment_variables "$1"
  if [ -n "${src_net:-}" ]; then
    add_exclusion_rule_to_iptables "$3"
  else
    while read -r all_nets; do
      src_net="$all_nets" add_exclusion_rule_to_iptables "$3"
    done <<< "$(xargs -n1 <<< "${2}")"
  fi
)
# END DMZ exclusion functions
trap 'rm -f /tmp/wg_started' EXIT
#
# Create wireguard DMZ
#
iptables --new WG_DMZ
iptables --new WG_DMZ_allow
# DMZ config with rfc1918 private IPv4 spaces
iptables -I WG_DMZ -j WG_DMZ_allow
iptables -A WG_DMZ -d 10.0.0.0/8 -j REJECT --reject-with icmp-host-prohibited
iptables -A WG_DMZ -d 172.16.0.0/12 -j REJECT --reject-with icmp-host-prohibited
iptables -A WG_DMZ -d 192.168.0.0/16 -j REJECT --reject-with icmp-host-prohibited
iptables -A WG_DMZ -j RETURN
#
# Watch for changes
#
export old_interface_checksum=0
while true; do
  if [ ! -f /tmp/wg_started ]; then
    until [ -f /etc/wireguard/wg0.conf ]; do
      sleep 1
    done
    touch /tmp/wg_started
    refresh_wireguard
  elif [ -f /etc/wireguard/wg0.conf ] && inotifywait -e modify /etc/wireguard/wg0.conf; then
    refresh_wireguard
  fi
  sleep 1;
done
